# IoTProject/setup/gateway.service
[Unit]
Description=IoT Mesh Gateway
After=network.target bluetooth.target mosquitto.service
Wants=mosquitto.service bluetooth.service

[Service]
Type=simple
User=__USER__
WorkingDirectory=__PROJECT_DIR__
# Add --zone "Zone Name" to tag messages with this gateway's location
ExecStart=__PROJECT_DIR__/venv/bin/python3 -m gateway.main --bluetooth --mqtt --zone "__ZONE__"
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
